# controller_rpi/systemd/controller-ai.service
[Unit]
Description=Controller AI Service (Raspberry Pi)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/MicroVIT/controller_rpi/src
# Use Orin's Ollama (already running for Robot1). Set OLLAMA_HOST in config/.env to override.
Environment="OLLAMA_HOST=http://10.13.68.159:11434"
Environment="OLLAMA_MODEL=qwen2.5:0.5b"
Environment="AI_TIMEOUT=60"
ExecStart=/bin/bash -c 'source /home/pi/MicroVIT/controller_rpi/venv/bin/activate && [ -f ../config/.env ] && set -a && . ../config/.env && set +a; exec python3 controller_ai_service.py'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
